# home/.config/systemd/user/taskwarrior_sync.service
[Unit]
Description=sync taskwarrior to server

[Service]
Type=oneshot
ExecStart=task sync

